# opt/adserver/adserver-admin.service
[Unit]
Description=HLS Ad Server Admin WebUI
After=network.target adserver.service
StartLimitIntervalSec=0

[Service]
User=media
Group=media
WorkingDirectory=/opt/adserver
Environment="PATH=/opt/adserver/venv/bin"
EnvironmentFile=/opt/adserver/.env
ExecStart=/opt/adserver/venv/bin/python3 -m uvicorn admin_app:app --host ${AD_ADMIN_HOST} --port ${AD_ADMIN_PORT} --workers 1
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
