# etc/systemd/system/adserver-admin.service
[Unit]
Description=AdServer Admin Web Interface
After=network.target adserver.service

[Service]
User=media
Group=media
WorkingDirectory=/opt/adserver
Environment="PATH=/opt/adserver/venv/bin"
ExecStart=/usr/bin/python3 -m uvicorn admin_app:app --host 127.0.0.1 --port 8089 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
